# deploy/resy-bot.service
[Unit]
Description=Resy bot (drop sniper + watchlist poller)
After=network-online.target
Wants=network-online.target
# Don't thrash if it crash-loops (e.g. bad creds).
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=resy
WorkingDirectory=/opt/resy-bot
# Pings a Slack/Discord webhook on a successful booking (optional).
Environment=RESY_BOT_WEBHOOK=
ExecStart=/opt/resy-bot/.venv/bin/python -m resy_bot run --mode both
Restart=always
RestartSec=5

# Light hardening.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
